# scripts/backup/crontab
# Nightly at 02:00 UTC (docs/10-operations.md section 1). Output goes to the
# container's stdout/stderr (PID 1's fds) so `docker compose logs backup`
# shows dump/upload/prune progress and failures.
0 2 * * * /backup.sh >>/proc/1/fd/1 2>>/proc/1/fd/2
